# configs/power/systemd-system/g14-gpu-boot-mode.service
[Unit]
Description=G14 GPU boot mode selector (AC/DC)
DefaultDependencies=no
Before=supergfxd.service
Before=display-manager.service
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/g14-gpu-boot-mode.sh
TimeoutStartSec=10

[Install]
WantedBy=supergfxd.service
